# .crontab
#Cron tasks for sitegen updates
* * * * * bash /sitegen/utils/siteup.sh /sitegen/utils > /sitegen/utils/cron.log
